# upstart/cerberus.styx.conf
description "cerberus' daemon"

start on runlevel 2
stop on runlevel [016]

respawn

exec sudo -u ubuntu styx.python /var/www/cerberus.fcgi
