# upstart/live.monet.conf
description "live's daemon"

start on starting client.monet
stop on stopped client.monet

respawn
respawn limit 10 5

exec sudo node /home/ubuntu/monet/live/server.js >> /var/www/log/live.log
