# upstart/live.memoize.conf
description "live's daemon"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 10 5

exec sudo node /home/ubuntu/memoize/live/server.js >> /var/www/log/live.log
