# upstart/mongodb.styx.conf
description "mongodb daemon"

start on runlevel 2
stop on runlevel [016]

respawn

exec mongod
