# upstart/nginx.monet.conf
description "nginx daemon"

start on starting client.monet
stop on stopped client.monet

respawn

exec /usr/sbin/nginx -c /etc/nginx/nginx.conf  -g "daemon off;"
