# upstart/nginx.memoize.conf
description "nginx daemon"

start on runlevel 2
stop on runlevel [016]

respawn

exec /usr/sbin/nginx -c /etc/nginx/nginx.conf  -g "daemon off;"
